# ha-vrrp/files/etc/hotplug.d/iface/95-ha-vrrp-apply
#!/bin/sh
get_dev() {
  local sec="$1"
  local IFACE USE_VLAN VLAN_ID DEV
  IFACE="$(uci -q get ha_vrrp.$sec.iface || echo)"
  [ -z "$IFACE" ] && return 1
  USE_VLAN="$(uci -q get ha_vrrp.$sec.use_vlan || echo 0)"
  VLAN_ID="$(uci -q get ha_vrrp.$sec.vlan_id || echo)"
  DEV="$IFACE"
  [ "$USE_VLAN" = "1" ] && [ -n "$VLAN_ID" ] && DEV="${IFACE}.${VLAN_ID}"
  echo "$DEV"
}
MON=""
INSTANCES="$(uci -q show ha_vrrp | awk -F. '/^ha_vrrp\.[^=]+=/ {print $2}' | cut -d= -f1 | sort -u | grep -E '^inst_' || true)"
if [ -z "$INSTANCES" ]; then
  MON="$(get_dev core || true)"
else
  for s in $INSTANCES; do
    d="$(get_dev "$s" || true)"; [ -n "$d" ] && MON="$MON $d"
  done
fi
for d in $MON; do
  [ "$INTERFACE" = "$d" ] || continue
  /usr/sbin/ha-vrrp-apply && /etc/init.d/keepalived restart
  exit 0
done
