# usr/lib/ha-vrrp/scripts/ha-vrrp/files/etc/init.d/ha-vrrp
#!/bin/sh /etc/rc.common
START=69
STOP=11
USE_PROCD=1

start_service() {
  local enabled="$(uci -q get ha_vrrp.core.enabled || echo 1)"
  [ "$enabled" = "1" ] || { logger -t ha-vrrp "disabled"; return 0; }

  /usr/libexec/ha-vrrp/ensure_vlan.sh || true
  /usr/sbin/ha-vrrp-apply || true

  /etc/init.d/keepalived stop >/dev/null 2>&1 || true
  /etc/init.d/keepalived start

  /etc/init.d/ha-vrrp-syncd start >/dev/null 2>&1 || true
}

stop_service() {
  /etc/init.d/ha-vrrp-syncd stop >/dev/null 2>&1 || true
  /etc/init.d/keepalived stop >/dev/null 2>&1 || true
}
